# systemd/user/sgen_bot.timer
[Unit]
Description=Tweet a new startup to @sgen_startups every hour

[Timer]
OnCalendar=*:03:14

[Install]
WantedBy=timers.target
